# config/includes.chroot/etc/systemd/system/zfs-auto-install.service
[Unit]
Description=Automated ZFS-on-root installer (DESTRUCTIVE)
After=network-online.target local-fs.target
Wants=network-online.target
ConditionKernelCommandLine=zfs.auto=1

[Service]
Type=simple
Environment=DEBIAN_FRONTEND=noninteractive
ExecStart=/usr/local/sbin/install-zfs-root.sh
StandardOutput=journal+console
StandardError=journal+console
TTYPath=/dev/tty1
TTYReset=yes
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
